# conf/supervisor-app.ini
[program:app-uwsgi]
command = /usr/bin/uwsgi --ini /etc/uwsgi/sites/apiserver.ini

[program:nginx-app]
command = /usr/sbin/nginx